# contrib/gula-bot-restart.service
[Unit]
Description=Daily bot restart

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart --no-block gula-bot.service